# ops/ghostbus-publisher.service
[Unit]
Description=Ghost Bus dataset publish (build CSVs, push to the data repo)
# Ordering only, not a dependency: if both are queued in the same boot the
# classifier finishes first, so the day being published is fully classified.
After=ghostbus-classifier.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# /opt/ghost-bus is chowned to ubuntu at install time (RUNBOOK 2), so running
# as ubuntu rather than the systemd default of root does two things at once:
# it makes the checkout's owner match the process's UID, so git's "dubious
# ownership" check never fires in the first place (no safe.directory
# workaround needed) - and it shrinks the blast radius of the one credential
# on this VM that can write to GitHub. EnvironmentFile is still read by
# systemd itself before the process starts, so /etc/ghostbus.env staying
# root:root mode 600 is unaffected.
#
# ghostbus-poller.service and ghostbus-classifier.service also run as
# ubuntu, deliberately, not just this unit: a live test showed root-owned
# WAL/SHM files were in fact readable by an ubuntu reader, but an
# unexplained empirical result about mixed-UID access to a shared SQLite
# database is a bad foundation for the one component holding a repo-write
# credential - all three sharing one owner removes the question rather than
# resting on it.
User=ubuntu
WorkingDirectory=/opt/ghost-bus
EnvironmentFile=/etc/ghostbus.env
ExecStart=/opt/ghost-bus/ops/publish.sh
